# shared/menu/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin __USERNAME__ --noclear %I 38400 linux
Type=idle
